# ansible/test/expected/etc/cron.d/risno
nsible: fetch pubs from sites
1 15 * * * risno node /opt/risno/slurp/fetch_sites >>/var/opt/risno/fetch.log 2>&1
#Ansible: update locations
10 17 * * * risno python /opt/risno/update_locations.py >>/var/opt/risno/upd_locations.log 2>&1
#Ansible: update types
1 17 * * * risno python /opt/risno/update_type.py >>/var/opt/risno/upd_types.log 2>&1
#Ansible: expire pubs
1 1 * * * risno python /opt/risno/expire_pubs.py --logic-immo --se-loger --paru-vendu --avendre-alouer --pages-jaunes --immo-street --belle-immobilier >>/var/opt/risno/expire.log 2>&1
#Ansible: purge expired pubs
1 14 * * * risno python /opt/risno/purge_pubs.py >>/var/opt/risno/purge.log 2>&1

